# packages/maintenance/skel/etc/omd/service_autorestart.cron
#
# this file will be managed by the SERVICE_AUTORESTART hook
#
# restart omd services if needed
* * * * *    /usr/bin/env lib/monitoring-plugins/check_omd.py --heal | perl -pe 'printf("[\%s] ", scalar localtime)' >> var/log/service_autorestart.log 2>&1
