# packages/omd/skel/etc/cron.d/clean_core_files
# remove core files older than 7 days
0 0 * * * find $OMD_ROOT/ -maxdepth 1 -type f -name core\* -mtime +7 -delete
