# packages/snmptrapd/samplicate.service
[Unit]
Description=Samplicate snmptrapd forwarder
After=samplicate_watch
Requires=samplicate_watch

[Service]
EnvironmentFile=/tmp/samplicate_targets.env
ExecStart=/opt/omd/versions/default/bin/samplicate -m /var/run/samplicate.pid -S -p 162 -f $SAMPLICATE_TARGETS
Type=forking
PIDFile=/var/run/samplicate.pid

[Install]
WantedBy=multi-user.target
