# packages/victoriametrics/skel/etc/cron.d/victoriametrics
10 * * * * ${OMD_ROOT}/bin/vmbackup-cron
05 * * * * ${OMD_ROOT}/bin/vmbackup-cleanup-cron

