# packages/maintenance/skel/etc/cron.d/logrotate
#
# Daily Logrotate
#
0 0 * * * [ -x /usr/sbin/logrotate ] && /usr/sbin/logrotate -s ###ROOT###/tmp/logrotate.state ###ROOT###/etc/logrotate.conf >/dev/null 2>&1
