# packages/sakuli/setup/omd/etc/cron.d/sakuli_screenshot_cleanup
MAILTO=""
0 1 * * * $OMD_ROOT/local/bin/sakuli_screenshot_cleanup.sh
